# image_files/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=*-*-* 04:00:00
